# etc/linguaevalia.service
[Unit]
Description=Lingua evalia

[Service]
User=nobody
Group=nobody
ExecStart=/home/lu/k/linguaevalia/linguaevalia server -p 8042
WorkingDirectory=/home/lu/k/linguaevalia

[Install]
WantedBy=multi-user.target
